# iso/lykon/airootfs/etc/systemd/system/plymouth-start.service.d/override.conf
[Service]
ExecStartPre=/usr/bin/plymouth --wait
StandardOutput=null
StandardError=null
